# templates/com.pi.memory.consolidate.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Managed automatically by pi-memory (`scheduler sync` on postinstall, init, session_start).
  Reference copy — normally do not install by hand.

  Manual test after sync:
    launchctl start com.pi.memory.maintenance

  Logs: __HOME__/.pi/pi-memory-data/logs/maintenance.{log,err.log}
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pi.memory.maintenance</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>mkdir -p '__HOME__/.pi/pi-memory-data/logs' &amp;&amp; exec '__NODE__' '__CLI__' maintenance --cron --verbose >> '__HOME__/.pi/pi-memory-data/logs/maintenance.log' 2>> '__HOME__/.pi/pi-memory-data/logs/maintenance.err.log'</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PI_MEMORY_ENV_FILE</key>
    <string>__HOME__/.pi/agent/pi-memory.env</string>
    <key>PI_MEMORY_AGENT_DIR</key>
    <string>__HOME__/.pi/pi-memory-data</string>
  </dict>
</dict>
</plist>
